# .docker/develop/php/supervisord.conf
[supervisord]
nodaemon=true

[program:octane]
command=php -d variables_order=EGPCS /var/www/artisan octane:start --server=roadrunner --host=0.0.0.0 --port=8080
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
